# onvif-voice-assistant/rootfs/etc/services.d/onvif-voice-assistant/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup when service stops
# ==============================================================================

bashio::log.warning "ONVIF Voice Assistant stopped!"
