# onvif-voice-assistant/rootfs/etc/services.d/onvif-voice-assistant/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start ONVIF Voice Assistant service
# ==============================================================================

bashio::log.info "Starting ONVIF Voice Assistant..."

# Executar aplicação Python
exec python3 /app/app.py
